# extensions/profile-avatars/sidecar/profile-avatars-sidecar.service
[Unit]
# Profile Avatars sidecar — loopback backend for the profile-avatars WebUI extension.
# The WebUI proxies /api/extensions/profile-avatars/sidecar/* here after consent.
# Install: cp to ~/.config/systemd/user/ then `systemctl --user enable --now profile-avatars-sidecar`
Description=Profile Avatars sidecar (per-profile avatar backend for the profile-avatars WebUI extension)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 %h/.hermes/webui/extensions/profile-avatars/sidecar/sidecar.py
WorkingDirectory=%h/.hermes/webui/extensions/profile-avatars/sidecar
# Share the WebUI's state dir so existing avatars.db (the 7 uploaded avatars) carries over.
Environment=HERMES_AVATARS_STATE_DIR=%h/.hermes/webui
Environment=HERMES_AVATARS_SIDECAR_PORT=17798
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
