# extensions/jarvis-voice/sidecar/jarvis-voice-sidecar.service
[Unit]
Description=Jarvis Voice Hermes WebUI sidecar
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/.hermes/webui/extensions/jarvis-voice/sidecar
Environment=HERMES_WEBUI_STATE_DIR=%h/.hermes/webui
ExecStart=/usr/bin/python3 -S -u sidecar.py
Restart=on-failure

[Install]
WantedBy=default.target
